# jstat.cron
# cron jstat collect once per minute
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin:/usr/java/latest/bin
* * * * * root /usr/local/nagios/bin/send_jstat.sh > /dev/null 2>&1
